# debian/postinst
#!/bin/sh
#

update-rc.d ftpproxy start 22 2 3 4 5 . stop 22 0 1 6 .
exit 0

